# manifests/upstart.conf.erb
description "Gitlab-Pivotal bridge"
author      "Ilkka Laukkanen <ilkka.s.laukkanen@gmail.com>"

start on    (local-filesystems and net-device-up IFACE=eth0)
stop on     shutdown

respawn
respawn limit 5 60

script
  exec sudo -u gitlab-bridge <%= @app_path %>/bin/gitlabpivotalbridge -DapplyEvolutions.default=true -Dhttp.port=<%= @play_port %>
end script

